# scripts/service/couchpotato/run
#!/bin/sh


exec /sbin/setuser couchpotato /usr/bin/python2.7 /couchpotato/CouchPotato.py  --config_file=/config/config.ini --data_dir=/config/data

